# qlogic-qla2xxx-10.02.14.01_k/extras/qla2xxx-nvmefc-boot-connection.service
[Unit]
Description=Auto-connect to sybsystems on FC-NVMe devices during Boot

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/sbin/qla2xxx/qla2xxx_nvme_boot_connections.sh"

[Install]
WantedBy=default.target
